# conf/svalinn.service
# -*- mode: conf -*-

[Unit]
Description=The receiver of data from XMiDT Caduceus.
Documentation=https://github.com/xmidt-org/svalinn
After=network.target

[Service]
EnvironmentFile=-/etc/default/svalinn
User=svalinn
ExecStart=/usr/bin/svalinn \
          $SVALINN_OPTS
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
LimitNOFILE=infinity

[Install]
WantedBy=multi-user.target
